# examples/linux/systemd/kardbrd-agent.service
[Unit]
Description=Kardbrd Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/kardbrd-agent/agent.env
ExecStartPre=-/bin/sh -c '${CONTAINER_CMD:-podman} rm -f kardbrd-agent'
ExecStart=/bin/sh -c '${CONTAINER_CMD:-podman} run --rm --name kardbrd-agent \
    -e KARDBRD_ID \
    -e KARDBRD_TOKEN \
    -e KARDBRD_AGENT \
    -e KARDBRD_URL \
    -e ANTHROPIC_API_KEY \
    -e LOG_LEVEL \
    -e GIT_AUTHOR_NAME \
    -e GIT_AUTHOR_EMAIL \
    -e GIT_COMMITTER_NAME \
    -e GIT_COMMITTER_EMAIL \
    -v %h/.local/share/kardbrd-agent/workspaces:/home/agent/workspaces \
    -v %h/.claude:/home/agent/.claude \
    -v %h/.ssh/kardbrd-agent:/home/agent/.ssh/id_ed25519:ro \
    -v %h/.ssh/kardbrd-agent-config:/home/agent/.ssh/config:ro \
    kardbrd-agent start \
        --cwd /home/agent/workspaces/repo'
ExecStop=/bin/sh -c '${CONTAINER_CMD:-podman} stop -t 600 kardbrd-agent'
TimeoutStopSec=610
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
