# examples/linux/systemd/kardbrd-agent-prune.timer
[Unit]
Description=Weekly prune of unused container images

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
