# examples/linux/systemd/kardbrd-agent-update.service
[Unit]
Description=Update kardbrd-agent container image

[Service]
Type=oneshot
ExecStart=%h/.local/bin/kardbrd-agent-update.sh
